# openwrt/files/etc/init.d/metaclash
#!/bin/sh /etc/rc.common
# /etc/init.d/clashforge — ClashForge procd service script

USE_PROCD=1
START=99
STOP=10

PROG=/usr/bin/clashforge
CONFIG=/etc/metaclash/config.toml
PIDFILE=/var/run/metaclash/metaclash.pid

# Bundled mihomo (installed with this package)
BUNDLED_MIHOMO=/usr/bin/mihomo-clashforge

find_mihomo() {
    # Priority 1: System mihomo already installed (user preference)
    if [ -x "/usr/bin/mihomo" ] && /usr/bin/mihomo -v > /dev/null 2>&1; then
        echo "/usr/bin/mihomo"
        return
    fi
    # Priority 2: OpenClash's mihomo (coexistence mode)
    if [ -x "/usr/bin/clash" ] && /usr/bin/clash -v > /dev/null 2>&1; then
        echo "/usr/bin/clash"
        return
    fi
    # Priority 3: Bundled mihomo (shipped with this IPK)
    if [ -x "$BUNDLED_MIHOMO" ] && $BUNDLED_MIHOMO -v > /dev/null 2>&1; then
        echo "$BUNDLED_MIHOMO"
        return
    fi
    return 1
}

patch_config_mihomo() {
    local mihomo_bin="$1"
    # Update the core.binary path in config so clashforge uses the right mihomo
    if [ -f "$CONFIG" ] && command -v sed > /dev/null 2>&1; then
        sed -i "s|^binary = .*|binary = \"$mihomo_bin\"|" "$CONFIG" 2>/dev/null || true
    fi
}

start_service() {
    local mihomo_bin
    mihomo_bin=$(find_mihomo)

    if [ -z "$mihomo_bin" ]; then
        echo "ClashForge: No mihomo binary found. Check /etc/metaclash/config.toml core.binary."
    else
        echo "ClashForge: Using mihomo at $mihomo_bin"
        patch_config_mihomo "$mihomo_bin"
    fi

    mkdir -p /var/run/metaclash

    procd_open_instance
    procd_set_param command "$PROG" -config "$CONFIG"
    procd_set_param pidfile "$PIDFILE"
    procd_set_param respawn 3600 5 3
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
}

reload_service() {
    local pid
    pid=$(cat "$PIDFILE" 2>/dev/null)
    [ -n "$pid" ] && kill -HUP "$pid" 2>/dev/null && echo "ClashForge: reloaded"
}

service_triggers() {
    procd_add_reload_trigger "network"
}

boot() {
    start "$@"
}
