# openwrt/files/etc/init.d/metaclash
#!/bin/sh /etc/rc.common
# /etc/init.d/clashforge — ClashForge procd service script
# Designed to coexist safely with OpenClash

USE_PROCD=1
START=99
STOP=10

PROG=/usr/bin/clashforge
CONFIG=/etc/metaclash/config.toml
PIDFILE=/var/run/metaclash/metaclash.pid
LOGFILE=/var/log/clashforge.log
MAX_LOG_KB=512
MEM_FLOOR_KB=32768

log() {
    logger -t clashforge "$*"
}

find_mihomo() {
    # Only check if executable exists — do NOT run -v (exit code is non-zero)
    # Priority 1: User-installed system mihomo
    if [ -x "/usr/bin/mihomo" ]; then
        echo "/usr/bin/mihomo"; return 0
    fi
    # Priority 2: Bundled mihomo (shipped with this IPK)
    if [ -x "/usr/bin/mihomo-clashforge" ]; then
        echo "/usr/bin/mihomo-clashforge"; return 0
    fi
    return 1
}

check_memory() {
    local free_kb
    free_kb=$(awk '/MemAvailable/{print $2}' /proc/meminfo 2>/dev/null || \
              awk '/MemFree/{print $2}' /proc/meminfo 2>/dev/null || echo 999999)
    if [ "$free_kb" -lt "$MEM_FLOOR_KB" ]; then
        log "ERROR: insufficient memory (${free_kb}KB free). Aborting start."
        return 1
    fi
    return 0
}

patch_config_mihomo() {
    local bin="$1"
    [ -f "$CONFIG" ] || return
    sed -i "s|^binary = .*|binary = \"$bin\"|" "$CONFIG" 2>/dev/null || true
}

coexist_with_openclash() {
    if /etc/init.d/openclash status 2>/dev/null | grep -q "running"; then
        log "OpenClash detected — switching to coexistence port set"
        if grep -q "^http = 7890" "$CONFIG" 2>/dev/null; then
            sed -i \
                -e 's/^http = 7890$/http = 17890/' \
                -e 's/^socks = 7891$/socks = 17891/' \
                -e 's/^mixed = 7893$/mixed = 17893/' \
                -e 's/^redir = 7892$/redir = 17892/' \
                -e 's/^tproxy = 7895$/tproxy = 17895/' \
                "$CONFIG" 2>/dev/null || true
        fi
        if grep -q '^mode = "tproxy"' "$CONFIG" 2>/dev/null; then
            sed -i 's/^mode = "tproxy"$/mode = "none"/' "$CONFIG" 2>/dev/null || true
        fi
    fi
}

start_service() {
    check_memory || return 1

    mkdir -p /var/run/metaclash
    mkdir -p /etc/metaclash

    # Clean up stale pidfile
    rm -f "$PIDFILE"

    [ -f "$CONFIG" ] || \
        touch "$CONFIG" 2>/dev/null || true
    [ -f /etc/metaclash/overrides.yaml ] || \
        touch /etc/metaclash/overrides.yaml 2>/dev/null || true

    local mihomo_bin
    mihomo_bin=$(find_mihomo)
    if [ -z "$mihomo_bin" ]; then
        log "ERROR: no mihomo binary found at /usr/bin/mihomo or /usr/bin/mihomo-clashforge"
        return 1
    fi
    log "Using mihomo: $mihomo_bin"
    patch_config_mihomo "$mihomo_bin"

    coexist_with_openclash

    procd_open_instance
    procd_set_param command "$PROG" -config "$CONFIG"
    procd_set_param pidfile "$PIDFILE"
    procd_set_param respawn 3600 10 3
    # Don't pipe stdout/stderr to console — clashforge writes its own log
    procd_close_instance
    log "clashforge started"
}

stop_service() {
    log "Stopping clashforge..."
}

reload_service() {
    local pid
    pid=$(cat "$PIDFILE" 2>/dev/null)
    if [ -n "$pid" ] && kill -0 "$pid" 2>/dev/null; then
        kill -HUP "$pid"
        log "Reloaded (SIGHUP sent to pid $pid)"
    else
        log "Not running, starting instead"
        start
    fi
}

service_triggers() {
    procd_add_reload_trigger "network"
}

boot() {
    sleep 5
    start
}
